# OpenBB/conda/pkgs/conda-26.1.1-py313h8f79df9_0/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/runner/miniforge3/conda-bld/conda_1772191422035/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pla/bin/conda";
setenv _CONDA_EXE "/Users/runner/miniforge3/conda-bld/conda_1772191422035/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pla/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/runner/miniforge3/conda-bld/conda_1772191422035/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pla/bin/python";
setenv _CONDA_ROOT "/Users/runner/miniforge3/conda-bld/conda_1772191422035/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pla";

source "/Users/runner/miniforge3/conda-bld/conda_1772191422035/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_pla/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"